# pkg/systemd/gddns.service
[Unit]
Description=gddns dynamic DNS updater

[Service]
User=gddns
Group=gddns
ExecStart=/usr/bin/gddns daemon

[Install]
WantedBy=multi-user.target
